# Jaar_2/Automation/pe/aut2223-pe1-RasmusLesebergPXL/app.service.j2
[Unit]
Description=app.py service

[Service]
User=vagrant
WorkingDirectory=/opt/pyramid
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target